# services/ai_api/com.astock.ai-api.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchAgent 模板 —— ai_api FastAPI 服务常驻.

  使用方法 (由 team-lead 配):
    1. 把 {{PROJECT_ROOT}} 替换为绝对路径 (e.g. /Users/wujiangjingcai/claude code/量化)
    2. 把 {{HOME}} 替换为 $HOME
    3. 确认 `codex login status` 已登录（服务复用本机登录态）
    4. 复制到 ~/Library/LaunchAgents/com.astock.ai-api.plist
    5. launchctl load ~/Library/LaunchAgents/com.astock.ai-api.plist
    6. launchctl start com.astock.ai-api

  日志:
    StandardOutPath / StandardErrorPath 都落到项目 .omc/logs/ 下
    日切 / size cap 由 launchd 自动管 (没配 LimitLoadTo, 走默认)

  失败重启:
    KeepAlive=true + Crashed=true → 进程崩了自动拉起
    ThrottleInterval=10 → 防爆: 崩了至少等 10s 再重启
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.astock.ai-api</string>

  <key>WorkingDirectory</key>
  <string>{{PROJECT_ROOT}}</string>

  <key>ProgramArguments</key>
  <array>
    <string>{{HOME}}/.local/bin/uv</string>
    <string>run</string>
    <string>uvicorn</string>
    <string>services.ai_api.app:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8000</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>{{HOME}}/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>LLM_PROVIDER</key>
    <string>codex</string>
    <key>CODEX_IGNORE_USER_CONFIG</key>
    <string>true</string>
    <key>CODEX_PROJECT_DOC_MAX_BYTES</key>
    <string>0</string>
    <key>CODEX_DISABLE_LOCAL_SKILLS</key>
    <string>true</string>
    <key>CODEX_ISOLATE_WORKDIR</key>
    <string>true</string>
    <key>CODEX_REASONING_EFFORT</key>
    <string>medium</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>Crashed</key>
    <true/>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>{{PROJECT_ROOT}}/.omc/logs/ai_api.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>{{PROJECT_ROOT}}/.omc/logs/ai_api.stderr.log</string>
</dict>
</plist>
